# rootfs/etc/crontabs/nobody
30 6 * * *    /usr/sbin/logrotate -v /etc/logrotate.conf